# zeroclaw-bridge.service.template
[Unit]
Description=ZeroClaw HTTP Bridge for StackChan
After=network-online.target
Wants=network-online.target

[Service]
# EDIT: replace /root/... paths with wherever you cloned zeroclaw-bridge
# and installed the zeroclaw binary. These defaults assume running as root
# with cargo-installed zeroclaw under /root/.cargo.
Type=simple
User=root
WorkingDirectory=/root/zeroclaw-bridge
Environment=ZEROCLAW_BIN=/root/.cargo/bin/zeroclaw
# Voice-daemon LLM tier identities. Bridge uses these in multi-profile mode
# (set via VOICE_LOCAL_PROFILE_KEY / VOICE_CLOUD_PROFILE_KEY below) to flip
# `[providers].fallback` between local llama-swap and OpenRouter on smart_mode.
Environment=DOTTY_DEFAULT_MODEL=qwen3.5:4b
# 27B cold-load on llama-swap (UD-Q4_K_XL @ 96K) is ~41s. Bridge's
# voice think_hard tool dispatches a direct call to qwen3.6:27b — bump
# this if you push 27B context higher or run on slower hardware.
Environment=VOICE_THINKER_TIMEOUT=90
Environment=SMART_MODEL=anthropic/claude-sonnet-4-6
Environment=VOICE_LOCAL_PROFILE_KEY=custom:http://<XIAOZHI_HOST>:8080/v1
Environment=VOICE_CLOUD_PROFILE_KEY=custom:https://openrouter.ai/api/v1
# Mirror of xiaozhi-server's PERSONA env var so the dashboard can show what's
# active. Keep these two in sync; the actual persona load happens in the
# xiaozhi-server container (PERSONA env in docker-compose.yml).
Environment=DOTTY_ACTIVE_PERSONA=default
# Optional: context injection for voice turns (date/time is always injected)
# Environment=WEATHER_LOCATION=Brisbane
# Environment=TZ=Australia/Brisbane
# Environment=CALENDAR_ID=<CALENDAR_ID_1>,<CALENDAR_ID_2>
# Environment=CALENDAR_SA_PATH=/root/.zeroclaw/secrets/google-calendar-sa.json
ExecStart=/root/zeroclaw-bridge/.venv/bin/python /root/zeroclaw-bridge/bridge.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
